# desktop-sync/launchd/com.tgbg.desktop-sync.hourly.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
  <dict>
    <key>Label</key>
    <string>com.tgbg.desktop-sync.hourly</string>

    <key>WorkingDirectory</key>
    <string>/ABSOLUTE/PATH/TO/tgbg/desktop-sync</string>

    <key>ProgramArguments</key>
    <array>
      <string>/bin/bash</string>
      <string>/ABSOLUTE/PATH/TO/tgbg/desktop-sync/run_sync_if_online_launchd.sh</string>
    </array>

    <!-- Run every hour -->
    <key>StartInterval</key>
    <integer>3600</integer>

    <key>RunAtLoad</key>
    <true/>

    <key>StandardOutPath</key>
    <string>/ABSOLUTE/PATH/TO/Library/Logs/tgbg-sync/stdout.log</string>
    <key>StandardErrorPath</key>
    <string>/ABSOLUTE/PATH/TO/Library/Logs/tgbg-sync/stderr.log</string>
  </dict>
</plist>

